# setup/systemd/shallot-alert-assess.timer
[Unit]
Description=Run Security Shallots alert-quality assessment every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
AccuracySec=2min
Persistent=true
Unit=shallot-alert-assess.service

[Install]
WantedBy=timers.target
