# setup/systemd/shallot-watchdog.service
[Unit]
Description=Security Shallots HTTP liveness watchdog
After=shallotd-home.service

[Service]
Type=oneshot
User=om
ExecStart=/home/user/security-shallots/.venv/bin/python /home/user/security-shallots/tools/shallot_watchdog.py
StandardOutput=journal
StandardError=journal
