# deploy/systemd/shallot-posture-scan.timer
[Unit]
Description=Run Security Shallots posture scan every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=shallot-posture-scan.service

[Install]
WantedBy=timers.target

